# api/docker/cron/tasks.cron
*/1 * * * * php /var/www/api/bin/console app:project:status-change